# scripts/com.lemon.lemon-ai-center.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Lemon AI Center — macOS Launch Agent
  Installs to: ~/Library/LaunchAgents/com.lemon.lemon-ai-center.plist

  To install:
    cp scripts/com.lemon.lemon-ai-center.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.lemon.lemon-ai-center.plist

  To start/stop:
    launchctl start com.lemon.lemon-ai-center
    launchctl stop com.lemon.lemon-ai-center

  To check status:
    launchctl list | grep lemon-ai-center

  To view logs:
    tail -f /tmp/lemon-ai-center.log
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lemon.lemon-ai-center</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/quantumcode/CODE/lemon-ai-center/start-prod.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/quantumcode/CODE/lemon-ai-center</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>/Users/quantumcode</string>
  </dict>

  <!-- Start automatically on login -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Restart automatically if it crashes -->
  <key>KeepAlive</key>
  <true/>

  <!-- Wait 5s before restart -->
  <key>ThrottleInterval</key>
  <integer>5</integer>

  <!-- Log output -->
  <key>StandardOutPath</key>
  <string>/tmp/lemon-ai-center.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/lemon-ai-center.error.log</string>
</dict>
</plist>
